# users/feralfile/systemd-services/chromium-ready.target
[Unit]
Description=Chromium successfully launched and cdp is up
Wants=feral-controld.service feral-setupd.service