# users/feralfile/systemd-services/system-ready.target
[Unit]
Description=Systemd successfully launched and ready
