# users/feralfile/.config/systemd/user/chromium-ready.target
[Unit]
Description=Chromium successfully launched and cdp is up
Wants=feral-connectd.service feral-setupd.service