# users/feralfile/.config/systemd/user/display-restore.service
[Unit]
Description=Display Restore Service - Auto-applies display rotation on HDMI events
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/home/feralfile/scripts/display-restore.sh
StandardOutput=append:/home/feralfile/.logs/display-restore.log
StandardError=append:/home/feralfile/.logs/display-restore.log
Restart=always
RestartSec=5

[Install]
WantedBy=default.target