# users/feralfile/.config/systemd/user/feral-watchdog.service
[Unit]
Description=FeralFile Kiosk Watchdog Daemon
After=feral-sys-monitord.service network.target
ConditionPathExists=/run/user/1000/bus

[Service]
Type=notify
ExecStart=/usr/bin/feral-watchdog --debug
Restart=always
RestartSec=5s
WatchdogSec=30s
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"
StandardOutput=append:/home/feralfile/.logs/watchdog.log
StandardError=append:/home/feralfile/.logs/watchdog.log

[Install]
WantedBy=default.target